# pi-companion/setup/vibe-sync-leds.service
[Unit]
Description=Vibe-Sync LED Controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Root is required for GPIO/PWM access on the Pi
User=root
WorkingDirectory=__PROJECT_DIR__
ExecStart=/usr/bin/python3 __PROJECT_DIR__/main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
